# test_data/launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.foobar.test</string>
  <key>Program</key>
  <string>/Test</string>
  <key>ProgramArguments</key>
  <array>
    <string>--flag</string>
    <string>arg1</string>
  </array>
  <key>StandardInPath</key>
  <string>/dev/null/in</string>
  <key>StandardOutPath</key>
  <string>/dev/null/out</string>
  <key>StandardErrorPath</key>
  <string>/dev/null/error</string>
  <key>WorkingDirectory</key>
  <string>/dev/null</string>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>3600</integer>
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Month</key>
      <integer>1</integer>
      <key>Hour</key>
      <integer>12</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Month</key>
      <integer>1</integer>
      <key>Hour</key>
      <integer>13</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
  </array>
  <key>StartOnMount</key>
  <true/>
  <key>WatchPaths</key>
  <array>
    <string>/tmp/dir</string>
    <string>/tmp/exist</string>
  </array>
  <key>QueueDirectories</key>
  <array>
    <string>/tmp/queue</string>
  </array>
  <key>KeepAlive</key>
  <true/>
  <key>UserName</key>
  <string>nobody</string>
  <key>GroupName</key>
  <string>nobody</string>
  <key>InitGroups</key>
  <true/>
  <key>Umask</key>
  <integer>23</integer>
  <key>RootDirectory</key>
  <string>/tmp/chroot</string>
  <key>AbandonProcessGroup</key>
  <true/>
  <key>ExitTimeOut</key>
  <integer>30</integer>
  <key>TimeOut</key>
  <integer>30</integer>
  <key>ThrottleInterval</key>
  <integer>30</integer>
  <key>Nice</key>
  <integer>-5</integer>
  <key>ProcessType</key>
  <string>Interactive</string>
</dict>
</plist>
